# syscleand.service
[Unit]
Description=SysClean Daemon
After=network.target

[Service]
Type=simple
ExecStart=/opt/sysclean/bin/syscleand
Restart=always
RestartSec=5
User=root
Group=root

# Hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/run /var/log /var/lib/sysclean /home/*/.local/share/sysclean

[Install]
WantedBy=multi-user.target
