# packaging/syscleand.service
[Unit]
Description=SysClean Execution Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/syscleand
Restart=on-failure
RestartSec=5
# Note: Root is required to execute apt, docker, and snap cleanups safely.
User=root
Group=root

# Hardened security features
ProtectSystem=full
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
